# content-manager.service
[Unit]
Description=Nova Labs Content Manager - Events & Reservations Display
After=network.target

[Service]
Type=notify
User=sklosky
Group=sklosky
WorkingDirectory=/home/sklosky/content-manager
Environment="PATH=/usr/bin:/usr/local/bin:/home/sklosky/.local/bin"
ExecStart=/usr/bin/python3 -m gunicorn --bind 0.0.0.0:8000 --workers 2 --timeout 30 --access-logfile - --error-logfile - wsgi:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=content-manager

[Install]
WantedBy=multi-user.target

